# services/add_devices.timer
[Unit]
Description=Run devices

[Timer]
OnBootSec=1s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target