# services/capture.timer
[Unit]
Description=Capture devices

[Timer]
OnBootSec=2s
# OnUnitActiveSec=5m

[Install]
WantedBy=timers.target